# casevault/supervisord.conf
[supervisord]
nodaemon=true

[program:web]
command=nginx -g 'daemon off;'

[program:app]
command=python run.py
directory=/app

[program:worker]
command=python -m worker
directory=/app
