# casevault/supervisord.conf
[supervisord]
nodaemon=true

[program:app]
command=python run.py
directory=/app

[program:worker]
command=python -m worker
directory=/app
